# systemd_units/rdkntp.target
[Unit]
Description=NTP Synchronization Complete
Documentation=man:systemd.special(7)
